# etl_project/main.cron
30 01 * * * python /home/itde1/WEXZ/main.py
